# target/linux/ramips/base-files/etc/uci-defaults/01_cs
#!/bin/sh
#
# Copyright (C) 2012 OpenWrt.org
#

. /lib/ramips.sh

board=$(ramips_board_name)

case "$board" in
ap0165)
	uci set wireless.radio0.disabled='0'	
	uci set wireless.default_radio0.network='lan'
	uci delete network.lan.ifname
	uci set network.lan._orig_ifname='wlan0'
	uci set network.lan._orig_bridge='true'
	uci delete network.lan_dev
	uci set wireless.default_radio0.ssid='MT7688'
	uci set wireless.default_radio0.encryption='psk2'
	uci set wireless.default_radio0.key='0123456789012'
	uci set wireless.default_radio0.wps_pushbutton='1'
	rm -rf /etc/rc.button/*
	touch /etc/config/button
	uci set button.wps=button
	uci set button.wps.button='reset'
	uci set button.wps.action='released'
	uci set button.wps.handler='echo "wps..." > /dev/console;cd /var/run/hostapd;for socket in *; do [ -S "$socket" ] || continue; hostapd_cli -i "$socket" wps_pbc; done;echo "1">/sys/class/leds/ap0165:green:wps/brightness; sleep 120; echo "0">/sys/class/leds/ap0165:green:wps/brightness'
	uci set button.wps.min='0'
	uci set button.wps.max='0'
	uci set button.reboot=button
	uci set button.reboot.button='reset'
	uci set button.reboot.action='released'
	uci set button.reboot.handler='echo "reboot" > /dev/console;sync;reboot'
	uci set button.reboot.min='2'
	uci set button.reboot.max='4'
	uci set button.factory_reset=button
	uci set button.factory_reset.button='reset'
	uci set button.factory_reset.action='released'
	uci set button.factory_reset.handler='echo "factory reset" > /dev/console;jffs2reset -y && reboot &'
	uci set button.factory_reset.min='5'
	uci set button.factory_reset.max='10'
	uci commit
	wifi down
	wifi up
	;;
esac
