# target/linux/ramips/base-files/etc/uci-defaults/01_cs
#!/bin/sh
#
#Copyright (C) 2018 OpenWrt.org
#

. /lib/ramips.sh

board=$(board_name)
case "$board" in
ip04296)
	uci set wireless.radio0.disable='0'
	uci set wireless.default_radio0.network='lan'
	uci set network.lan._orig_ifname='wlan0'
	uci set network.lan._orig_bridge='ture'
	uci set wireless.default_radio0.ssid='IP04296'
	uci set wireless.default_radio0.encryption='psk2'
	uci set wireless.default_radio0.key='0123456789'
	uci set wireless.default_radio0.wps_pushbutton='1'
	rm -rf /etc/rc.button/*
	touch /etc/config/button
	uci set button.wps=button
	uci set button.wps.button='wps'
	uci set button.wps.action='released'
	uci set button.wps.handler='echo "wps sucefull!!" >/dev/console;cd /var/run/hostpad;for socket in *;do [ -S "$socket" ] || continue;hostpad_cli -i "$socket" wps_pbc; done;echo "1">/sys/class/leds/IP04296:amber:status/brightness; sleep 120; echo "0">/sys/class/leds/IP04296:amber:status/brightness'
	uci set button.wps.min='0'
	uci set button.wps.max='2'
	uci set button.reboot=button
	uci set button.reboot.button='reset'
	uci set button.reboot.action='released'
	uci set button.reboot.handler='echo "reboot" > /dev/console;echo "timer">/sys/class/leds/IP04296:amber:status/trigger;echo "50">/sys/class/leds/IP04296:amber:status/delay_on;echo "50">/sys/class/leds/IP04296:amber:status/delay_off;sleep 3;echo "none">/sys/class/leds/IP04296:amber:status/trigger;sync;reboot'
	uci set button.reboot.min='0'
	uci set button.reboot.max='3'
	uci set button.factory_reset=button
	uci set button.factory_reset.button='reset'
	uci set button.factory_reset.action='released'
	uci set button.factory_reset.handler='echo "factory reset" > /dev/console;echo "timer">/sys/class/leds/IP04296:amber:status/trigger;echo "75">/sys/class/leds/IP04296:amber:status/delay_on;echo "75">/sys/class/leds/IP04296:amber:status/delay_off;sleep 6;echo "none">/sys/class/leds/IP04296:amber:status/trigger;jffs2reset -y && reboot &'
	uci set button.factory_reset.min='5'
	uci set button.factory_reset.max='10'
 	uci commit
	wifi down
	wifi up
        ;;
esac


